# run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
After=blockdev@dev-mmcblk1p4.target

[Mount]
Where=/
What=/dev/mmcblk1p4
Type=ext4
Options=defaults,noatime,errors=remount-ro
