# var/lib/dpkg/info/cgmanager.list
/etc
/etc/init.d
/etc/init.d/cgproxy
/etc/init.d/cgmanager
/etc/init
/etc/init/cgproxy.conf
/etc/init/cgmanager.conf
